# lightdm/dmlock.service
[Unit]
Description=DM Lock before sleep
Before=sleep.target

[Service]
ExecStart=/etc/lightdm/suslock.sh

[Install]
WantedBy=sleep.target
